# deploy/server_examples/crontab
Пример файла для создания событий по крону

*/5  *  *  *  * /usr/bin/php /usr/share/nginx/yii2-lottery/yii finance-ext/re-quene-payments
